# pool.service
# Canonical service definition for the Pool Controller.
# Deploy: after `git pull`, copy to /etc/systemd/system/pool.service,
# then: sudo systemctl daemon-reload && sudo systemctl restart pool.service
#
# Secrets (MQTT password) live in /etc/pool-controller/env — see README.md.

[Unit]
Description=Pool Controller
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/rpi-pool-controller/pool_controller.py
EnvironmentFile=/etc/pool-controller/env
Restart=always
RestartSec=10
User=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
